# packages/sysutils/mactool/system.d/mactool-bt@.service
[Unit]
Description=Configure Bluetooth MAC address for %I

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/mactool-config %I
